# app.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * . $HOME/.profile; cd /var/www && /usr/local/bin/php artisan schedule:run >> /tmp/cron.log 2>&1
